# tiktokbot.service
[Unit]
Description=TikTok Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/BotTikTok
ExecStart=/home/ubuntu/BotTikTok/venv/bin/python bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
